# ops/systemd/aiharu-daily-prompt.timer
[Unit]
Description=Run AIHaru local daily prompt at 08:30 Asia/Seoul

[Timer]
OnCalendar=*-*-* 08:30:00 Asia/Seoul
Persistent=true
AccuracySec=1min
Unit=aiharu-daily-prompt.service

[Install]
WantedBy=timers.target
