# ops/systemd/aiharu-collect-news.timer
[Unit]
Description=Run AIHaru local news collection at 12:30 Asia/Seoul

[Timer]
OnCalendar=*-*-* 12:30:00 Asia/Seoul
Persistent=true
AccuracySec=1min
Unit=aiharu-collect-news.service

[Install]
WantedBy=timers.target
