# ops/systemd/aiharu-collect-news.timer
[Unit]
Description=Run AIHaru news collection at 13:00 Asia/Seoul

[Timer]
OnCalendar=*-*-* 13:00:00 Asia/Seoul
Persistent=true
AccuracySec=1min
Unit=aiharu-collect-news.service

[Install]
WantedBy=timers.target
