# ops/systemd/aiharu-daily-prompt.timer
[Unit]
Description=Run AIHaru daily prompt at 09:00 Asia/Seoul

[Timer]
OnCalendar=*-*-* 09:00:00 Asia/Seoul
Persistent=true
AccuracySec=1min
Unit=aiharu-daily-prompt.service

[Install]
WantedBy=timers.target
